# cathedral-os/platforms/linux/systemd/cathedral-slow.socket
[Unit]
Description=Cathedral OS - Slow Brain Activation Socket

[Socket]
ListenStream=/run/cathedral/slow_brain.sock
